# deploy/tinkoff-bot.service
[Unit]
Description=Tinkoff trading bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Замени на свой путь к проекту, например /root/tinkoff_bot
WorkingDirectory=/root/tinkoff_bot
ExecStart=/root/tinkoff_bot/.venv/bin/python run_bot.py
Restart=always
RestartSec=10
Environment=PATH=/root/tinkoff_bot/.venv/bin:/usr/bin

[Install]
WantedBy=multi-user.target
